# root/etc/services.d/pgdvs/run
#!/usr/bin/execlineb -P
python /config/pgdvs/default.py /config/pgdvs/gdrive.db